# darwis_panel/docker_resources/supervisor/webserver.conf
[program:webserver]

command=/bin/bash -c "service php8.2-fpm start && service apache2 start"
redirect_stderr=true
