# code/darwis_panel/docker_resources/supervisor/pythonservice.conf
[program:pythonservice]
user=appuser 
command=/bin/bash -c "cd /app/dt/threat_intel/src/ && python3 service.py"
redirect_stderr=true
